# netfiremon-web.service
[Unit]
Description=Net-Fire-Monitor v3.9 – Web-Interface (unprivilegiert)
Documentation=https://github.com/netfiremon
After=network-online.target netfiremon.service
Wants=network-online.target
Requires=netfiremon.service

[Service]
Type=simple
User=netfiremon
Group=netfiremon
ExecStart=/opt/netfiremon/.venv/bin/gunicorn \
    "netfiremon_web:create_wsgi_app()" \
    --bind 0.0.0.0:5443 \
    --workers 2 \
    --timeout 30 \
    --certfile /opt/netfiremon/certs/cert.pem \
    --keyfile /opt/netfiremon/certs/key.pem \
    --access-logfile - \
    --error-logfile - \
    --log-level warning
WorkingDirectory=/opt/netfiremon
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=netfiremon-web

# Sicherheits-Härtung: Web-Prozess braucht keinerlei erhöhte Rechte.
# Alle Firewall-Aktionen werden über IPC (CMD_QUEUE_DIR) an den
# Monitor-Prozess (root) delegiert.
ProtectSystem=strict
ReadWritePaths=/opt/netfiremon/data
ProtectHome=true
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
