# netfiremon.service
[Unit]
Description=Net-Fire-Monitor v3.9 – Netzwerk-Monitor (root)
Documentation=https://github.com/netfiremon
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
ExecStart=/opt/netfiremon/.venv/bin/python /opt/netfiremon/netfiremon_terminal.py --auto
WorkingDirectory=/opt/netfiremon
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=netfiremon

# Sicherheits-Härtung
ProtectSystem=strict
ReadWritePaths=/opt/netfiremon/data /run
PrivateTmp=true
ProtectHome=true
NoNewPrivileges=false
# Root-Rechte sind erforderlich für:
#   - Scapy (Raw-Socket für Paket-Capture)
#   - iptables/ip6tables (Firewall-Regeln setzen/löschen)
# Der Web-Prozess läuft separat als unprivilegierter User.

# Capabilities statt vollem root (Principle of Least Privilege)
AmbientCapabilities=CAP_NET_RAW CAP_NET_ADMIN
CapabilityBoundingSet=CAP_NET_RAW CAP_NET_ADMIN

[Install]
WantedBy=multi-user.target
